# playground/iximiuz/image/use-practice-bootstrap.service
[Unit]
Description=Bootstrap USE Practice Lab
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/install -d -m 0755 /var/lib/use-practice
ExecStart=/usr/bin/chown -R laborant:laborant /opt/use-practice
ExecStart=/usr/bin/ln -sfn /opt/use-practice /home/laborant/use-practice
ExecStart=/usr/bin/chown -h laborant:laborant /home/laborant/use-practice
ExecStart=/usr/bin/test -x /usr/bin/iostat
ExecStart=/usr/bin/test -x /usr/bin/mpstat
ExecStart=/usr/bin/test -x /usr/bin/pidstat
ExecStart=/usr/bin/test -x /usr/bin/sar
ExecStart=/usr/bin/test -x /usr/local/bin/use-practice
ExecStart=/usr/bin/test -x /usr/local/bin/use-tool
ExecStart=/usr/bin/test -x /opt/use-practice/bin/uworker
ExecStart=/usr/bin/test -x /opt/use-practice/bin/updisk
ExecStart=/usr/bin/test -x /opt/use-practice/bin/uwait
ExecStart=/usr/bin/touch /var/lib/use-practice/ready
ExecStart=/usr/bin/chmod 0644 /var/lib/use-practice/ready
RemainAfterExit=yes
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
